# integrations/docker/images/stage-2/chip-build-linux-qemu/files/config/btvirt.service
[Unit]
Description=btvirt adapter
After=syslog.target network.target

[Service]
ExecStart=/usr/bin/btvilt -l2 -L
