# integrations/docker/images/stage-2/chip-build-linux-qemu/files/config/hostapd.service
[Unit]
After=wifi_master.service
Requires=wifi_master.service

[Service]
NetworkNamespacePath=/var/run/netns/tool
